# daemon/resources/razer.systemd.service.in
[Unit]
Description=Service to manage razer devices in userspace
Documentation=man:razer-service(8)

[Service]
Type=dbus
BusName=org.razer
ExecStart=##PREFIX##/bin/razer-service -Fv

[Install]
WantedBy=default.target
